# auto-prop-node/auto-propnode.timer
[Unit]
Description=Run Auto Propagation Node Selector hourly

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
